# backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

# ══════════════════════════════════════════════════════════════════════════════
# Xvfb Virtual Displays — one per carrier for isolated VNC tabs
# ══════════════════════════════════════════════════════════════════════════════

[program:xvfb_maersk]
command=Xvfb :99 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:xvfb_cma]
command=Xvfb :100 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:xvfb_one]
command=Xvfb :101 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:xvfb_hapag]
command=Xvfb :102 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# ══════════════════════════════════════════════════════════════════════════════
# x11vnc — VNC servers for each display
# ══════════════════════════════════════════════════════════════════════════════

[program:x11vnc_maersk]
command=x11vnc -display :99 -forever -shared -nopw -rfbport 5900 -xkb -noxrecord -noxfixes -noxdamage
autostart=true
autorestart=true
priority=20
startsecs=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:x11vnc_cma]
command=x11vnc -display :100 -forever -shared -nopw -rfbport 5901 -xkb -noxrecord -noxfixes -noxdamage
autostart=true
autorestart=true
priority=20
startsecs=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:x11vnc_one]
command=x11vnc -display :101 -forever -shared -nopw -rfbport 5902 -xkb -noxrecord -noxfixes -noxdamage
autostart=true
autorestart=true
priority=20
startsecs=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:x11vnc_hapag]
command=x11vnc -display :102 -forever -shared -nopw -rfbport 5903 -xkb -noxrecord -noxfixes -noxdamage
autostart=true
autorestart=true
priority=20
startsecs=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# ══════════════════════════════════════════════════════════════════════════════
# websockify — WebSocket bridges for each VNC server
# ══════════════════════════════════════════════════════════════════════════════

[program:websockify_maersk]
command=websockify --web /usr/share/novnc 6080 localhost:5900
autostart=true
autorestart=true
priority=30
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:websockify_cma]
command=websockify --web /usr/share/novnc 6081 localhost:5901
autostart=true
autorestart=true
priority=30
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:websockify_one]
command=websockify --web /usr/share/novnc 6082 localhost:5902
autostart=true
autorestart=true
priority=30
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:websockify_hapag]
command=websockify --web /usr/share/novnc 6083 localhost:5903
autostart=true
autorestart=true
priority=30
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# ══════════════════════════════════════════════════════════════════════════════
# nginx + FastAPI
# ══════════════════════════════════════════════════════════════════════════════

[program:nginx]
command=nginx -c /app/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:fastapi]
command=uvicorn main:app --host 127.0.0.1 --port 8001
directory=/app
autostart=true
autorestart=true
priority=50
environment=DISPLAY=":99"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
